# bm2015.conf
[program:worker]
directory=/home/pi/venus-lift-lights
command=python /home/pi/venus-lift-lights/worker.py
numprocs=1
redirect_stderr=true
stdout_logfile=/tmp/worker.log
user=pi
autorestart=true
startsecs=2
environment= PYTHONPATH="$PYTHONPATH:/usr/local/lib/python2.7/site-packages"
priority=3

[program:myo]
directory=/home/pi/venus-lift-lights
command=python /home/pi/venus-lift-lights/myo.py
numprocs=1
redirect_stderr=true
stdout_logfile=/tmp/myo.log
redirect_stderr=true
user=pi
autorestart=true
startsecs=10
environment= PYTHONPATH="/usr/local/lib/python2.7/site-packages:$PYTHONPATH"
priority=5

[program:beep]
command=python /home/pi/venus-lift-lights/beep.py
numprocs=1
redirect_stderr=true
stdout_logfile=/tmp/beep.log
user=pi
autorestart=true
startsecs=1
environment= PYTHONPATH="$PYTHONPATH:/usr/local/lib/python2.7/site-packages"
priority=2

[program:button]
command=python /home/pi/venus-lift-lights/button.py
numprocs=1
redirect_stderr=true
stdout_logfile=/tmp/button.log
user=pi
autorestart=true
startsecs=1
environment= PYTHONPATH="$PYTHONPATH:/usr/local/lib/python2.7/site-packages"
priority=10

